# recipes-bsp/tegra-binaries/tegra-redundant-boot/nv_update_verifier.service
[Unit]
Description=OTA bootloader update verifier

[Service]
Type=simple
ExecStart=/usr/sbin/nvbootctrl verify

[Install]
WantedBy=multi-user.target
