# recipes-bsp/tools/setup-nv-boot-control/esp.mount.in
[Unit]
Before=local-fs.target

[Mount]
What=/dev/disk/by-partlabel/esp
Where=@ESPMOUNT@
Type=vfat
