# recipes-bsp/tegra-binaries/tegra-nvsciipc/nv_nvsciipc_init.service
[Unit]
Description=nvsciipc initialization
After=sysinit.target
Before=multi-user.target

[Service]
ExecStart=/usr/bin/nvsciipc_init
Type=oneshot
