# external/virtualization-layer/recipes-containers/nvidia-container-toolkit/nvidia-container-toolkit/nvidia-container-setup.service.in
[Unit]
Description=Generate NVIDIA container runtime configuration
Before=docker.socket containerd.service

[Service]
Type=oneshot
RemainAfterExit=yes
RuntimeDirectory=nvidia-container-runtime
RuntimeDirectoryPreserve=yes
ExecStart=@LIBEXECDIR@/nvidia-container-runtime/generate-config
