# recipes-security/optee/optee-client-tegra/nv-tee-supplicant.service.in
[Unit]
Description=OP-TEE Client Supplicant

[Service]
Type=simple
ExecStart=@sbindir@/tee-supplicant
Restart=always

[Install]
WantedBy=basic.target
