# recipes-bsp/l4t-usb-device-mode/l4t-usb-device-mode/98-usb-gadget-tty.rules
SUBSYSTEM=="tty", KERNEL=="ttyGS[0-9]", ACTION=="add", RUN+="/bin/systemctl start serial-getty@%k.service"
SUBSYSTEM=="tty", KERNEL=="ttyGS[0-9]", ACTION=="remove", RUN+="/bin/systemctl stop serial-getty@%k.service"
